# roles/cleversite/codes/templates/cleversite-codes.service.j2
[Unit]
Description=cleversite-codes
After=network-online.target

[Service]
WorkingDirectory={{app_current_release_path}}
ExecStart=/usr/bin/node --expose_gc --nouse-idle-notification server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=cleversite-codes
User=apache
Group=apache
Environment='NODE_ENV=production'
## limit CPU and RAM quota for our service
#CPUAccounting=true
#CPUQuota=10%
#MemoryAccounting=true
#MemoryLimit=50M


[Install]
WantedBy=multi-user.target